# .config/autostart/polybar.desktop
[Desktop Entry]
Type=Application
Name=polybar
Exec=sh -c '~/.config/polybar/i3cfg.sh'
StartupNotify=false
Terminal=false
